# package/satfi/src/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#checkgprs.sh
rmmod /lib/modules/3.10.49/power_mode.ko
sleep 1
insmod /lib/modules/3.10.49/power_mode.ko

exit 0
